# deploy/blog_server.ini
[program:blog_http]        # blog http服务
directory=/root/src
command=gunicorn -w 4 main:app -b :5856 --timeout=600  ; 程序启动命令
autostart=true       ; 在supervisord启动的时候也自动启动
startsecs=10         ; 启动10秒后没有异常退出，就表示进程正常启动了，默认为1秒
autorestart=true     ; 程序退出后自动重启,可选值：[unexpected,true,false]，默认为unexpected，表示进程意外杀死后才重启
startretries=3       ; 启动失败自动重试次数，默认是3
stdout_logfile_maxbytes=50MB
stderr_logfile_maxbytes=50MB
stdout_logfile=/root/logs/blog_http.log
stderr_logfile=/root/logs/blog_http.log
stdout_logfile_backups=7
stderr_logfile_backups=7
